# crontabs
# 凌晨 4 点清除老旧版本的 OTA 包
0 4 * * * sh ~/MaaRelease/MaaAssistantArknights/clean_ota.sh
